# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:elastic]
user=nobody
command=/elasticsearch-1.0.1/bin/elasticsearch

